# root/etc/services.d/openpyn/run
#!/usr/bin/with-contenv bash

/openpyn/openpyn.exp
